# deploy/systemd/pulse-db-backup.service
[Unit]
Description=Pulse Pager Postgres backup (portable logical dump)
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
# Point this at wherever the repo is checked out on the host.
ExecStart=/opt/pulse-pager/deploy/db/backup.sh
# Where to write backups. Copy them off the Pi periodically; a backup on the same
# SD card as the data is not real safety.
Environment=PULSE_BACKUP_DIR=/var/backups/pulse
Environment=PULSE_BACKUP_KEEP=14
